# container/config/crontabs/csdha
* * * * * flock -n /tmp/laravel-queue.lock php -d max_execution_time=120 artisan queue:work --once -q 
* * * * * echo hello cron >> /tmp/cron_test.log

